# cos-uploader.service
[Unit]
Description=COS File Monitor & Uploader Service
Documentation=file:///opt/cos-uploader/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu

# 工作目录
WorkingDirectory=/opt/cos-uploader

# 启动命令
ExecStart=/opt/cos-uploader/cos-uploader -config /opt/cos-uploader/config.yaml

# 重启策略
Restart=always
RestartSec=10

# 进程超时（0 表示无限制）
TimeoutStartSec=0
TimeoutStopSec=30

# 资源限制
LimitNOFILE=65535
LimitNPROC=65535

# 日志和输出
StandardOutput=append:/var/log/cos-uploader/cos-uploader.log
StandardError=append:/var/log/cos-uploader/cos-uploader.error.log

# 环境变量
Environment="PATH=/opt/cos-uploader:/usr/local/bin:/usr/bin:/bin"

# 权限和安全相关
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=yes
# ReadWritePaths=/var/lib/cos-uploader /var/log/cos-uploader

# 自动重启失败的服务
SuccessExitStatus=0 143

# 进程追踪
OOMPolicy=continue

[Install]
WantedBy=multi-user.target
